# crontab
@reboot sleep 10 && $HOME/devops-academy/expose_docker.sh && sleep 10 && docker-compose -f $HOME/devops-academy/docker-compose.yml u
p -d
